# _templates/launchagent.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.agentic-os.scheduler</string>
    <key>ProgramArguments</key>
    <array>
        <string>{{node_bin}}</string>
        <string>{{repo_root}}/scripts/scheduler-tick.mjs</string>
    </array>
    <key>WorkingDirectory</key>
    <string>{{repo_root}}</string>
    <key>StartInterval</key>
    <integer>60</integer>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>{{repo_root}}/.claude/state/scheduler.out.log</string>
    <key>StandardErrorPath</key>
    <string>{{repo_root}}/.claude/state/scheduler.err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>{{path_env}}</string>
    </dict>
</dict>
</plist>
